# nix/packages/evo4/90-pulseaudio-evo4.rules
ACTION=="add", RUN+="RUN_TARGET" ENV{SYSTEMD_USER_WANTS}+="evo4.service"