# nix/packages/evo4/evo4.service
[Unit]
Description=Add evo4

[Service]
Type=simple
ExecStart=/bin/evo4
RemainAfterExit=no

[Install]
WantedBy=default.target